# files/systemd/kai-edge.service.tmpl
[Unit]
Description=kai edge runtime daemon
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=__KAI_USER__
Group=__KAI_GROUP__
SupplementaryGroups=audio
WorkingDirectory=__KAI_APP_DIR__
RuntimeDirectory=kai-edge
RuntimeDirectoryMode=0755
Environment=PYTHONUNBUFFERED=1
Environment=KAI_STATE_DIR=__KAI_STATE_DIR__
Environment=KAI_LOG_DIR=__KAI_LOG_DIR__
EnvironmentFile=-__KAI_CONFIG_DIR__/edge.env
ExecStart=__KAI_BIN_DIR__/kai-edge-daemon --env-file __KAI_CONFIG_DIR__/edge.env
Restart=on-failure
RestartSec=2
NoNewPrivileges=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
